# net-dns/pdnsd/files/pdnsd.online
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-dns/pdnsd/files/pdnsd.online,v 1.1.1.1 2005/11/30 09:50:16 chriswhite Exp $

depend() {
	need pdnsd net.${IFACE}
}

start() {
	ebegin "Starting pdnsd-online"
	/usr/sbin/pdnsd-ctl server all up
	eend $?
}

stop() {
	ebegin "Stopping pdnsd-online"
	/usr/sbin/pdnsd-ctl server all down
	eend $?
}
