# app-admin/hddtemp/files/hddtemp-init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-admin/hddtemp/files/hddtemp-init,v 1.1.1.1 2005/11/30 10:00:13 chriswhite Exp $

depend() {
	need localmount 
}

start() {
	ebegin "Starting hddtemp daemon"
	/sbin/start-stop-daemon --start --quiet --name hddtemp  \
	--exec ${HDDTEMP_EXEC} -- -d ${HDDTEMP_OPTS}

	eend $?
}

stop() {
	ebegin "Stopping hddtemp daemon"
	start-stop-daemon --stop --quiet --name hddtemp
	eend $?
}


