# net-misc/powerd/files/powerd-init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-misc/powerd/files/powerd-init,v 1.1.1.1 2005/11/30 09:55:03 chriswhite Exp $

start() {
	ebegin "Starting powerd"
	start-stop-daemon --start --exec /usr/sbin/powerd
	eend $?
}

stop() {
	ebegin "Stopping powerd"
	killall powerd
	eend $?
}
