# x11-misc/xhkeys/files/xhkeys
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/x11-misc/xhkeys/files/xhkeys,v 1.1.1.1 2005/11/30 09:40:24 chriswhite Exp $

start() {
	ebegin "Starting xhkeys"
	start-stop-daemon --start --quiet --exec /usr/bin/xhkeys	
	eend $?
}

stop() {
	ebegin "Stopping xhkeys"
	start-stop-daemon --stop --quiet --exec /usr/bin/xhkeys
	eend $?
}
