# net-misc/knock/files/knockd.initd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/net-misc/knock/files/knockd.initd,v 1.1.1.1 2005/11/30 09:55:39 chriswhite Exp $

depend() {
	need net
}

start() {
	ebegin "Starting knockd"
	start-stop-daemon --start --quiet --exec /usr/sbin/knockd -- ${OPTS}
	eend $?
}

stop() {
	ebegin "Stopping knockd"
	start-stop-daemon --stop --quiet --exec /usr/sbin/knockd
	eend $?
}
