# app-emulation/win4lin/files/win4lin.initd.new
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-emulation/win4lin/files/win4lin.initd.new,v 1.1.1.1 2005/11/30 10:08:56 chriswhite Exp $

depend() {
	need net
}

start() {
	ebegin "Starting Win4Lin"
	/etc/rc.d/init.d/Win4Lin start > /dev/null
	eend $?
}

stop() {
	ebegin "Stopping Win4Lin"
	/etc/rc.d/init.d/Win4Lin stop > /dev/null
	eend $?
}

