# sys-libs/libchipcard/files/chipcardd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-libs/libchipcard/files/chipcardd,v 1.1.1.1 2005/11/30 09:39:18 chriswhite Exp $

start() {
    ebegin "Starting libchipcard daemon"
    start-stop-daemon --start --quiet --exec /usr/sbin/chipcardd
    eend $?
}

stop() {
    ebegin "Stopping libchipcard daemon"
    start-stop-daemon --stop --quiet --exec /usr/sbin/chipcardd
    eend $?
}
